# user_salt/files/qubes-amqp.socket
[Unit]
Description=Minio TCP Socket

[Socket]
ListenStream=127.0.0.1:25672
Accept=true

[Install]
WantedBy=sockets.target

