# user_salt/files/qubes-amqp.service
[Unit]
Description=Minio Server Port

[Service]
ExecStart=qrexec-client-vm '' qubes.ConnectTCP+5672
StandardInput=socket
StandardOutput=inherit

