# user_salt/files/serviceunits/qubes-minio@.service
[Unit]
Description=Minio Server Port

[Service]
ExecStart=qrexec-client-vm '' qubes.ConnectTCP+9000
StandardInput=socket
StandardOutput=inherit

